# deploy/systemd/application-notifier-morning.service
[Unit]
Description=Application notifier morning reminder
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/application-notifier/application-notifier.env
ExecStart=/bin/sh -lc 'exec docker exec -e "CZM_BASE_URL=$CZM_BASE_URL" -e "CZM_API_KEY=$CZM_API_KEY" -e "CZM_TIMEZONE=$CZM_TIMEZONE" -e "CZM_CONFIG_PATH=$CZM_CONFIG_PATH" -e "OPENCLAW_BRIDGE_MODE=$OPENCLAW_BRIDGE_MODE" -e "OPENCLAW_BRIDGE_COMMAND=$OPENCLAW_BRIDGE_COMMAND" -e "OPENCLAW_BRIDGE_FALLBACK_COMMAND=$OPENCLAW_BRIDGE_FALLBACK_COMMAND" -e "OPENCLAW_BRIDGE_TARGET=$OPENCLAW_BRIDGE_TARGET" -e "OPENCLAW_BRIDGE_TIMEOUT_SECONDS=$OPENCLAW_BRIDGE_TIMEOUT_SECONDS" -e "OPENCLAW_REMINDER_COMPOSE_COMMAND=$OPENCLAW_REMINDER_COMPOSE_COMMAND" -e "OPENCLAW_PYTHON_BIN=$OPENCLAW_PYTHON_BIN" -e "APPLICATION_NOTIFIER_LOCK_PATH=$APPLICATION_NOTIFIER_LOCK_PATH" "${OPENCLAW_CONTAINER_NAME:-openclaw}" "${OPENCLAW_PYTHON_BIN:-python3}" "${NOTIFIER_DEPLOY_PATH_IN_CONTAINER:-/app/application-notifier}/run_reminder.py" --slot morning'
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
