# deploy/systemd/application-notifier-evening.timer
[Unit]
Description=Run application notifier evening reminder at 20:00

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true
Unit=application-notifier-evening.service

[Install]
WantedBy=timers.target

